# research_assistant/peak-research.service
[Unit]
Description=PEAK Research Assistant Service
After=network.target

[Service]
# Set the working directory where the script and input files are located
WorkingDirectory=/home/ubuntu/PEAK-Assistant/autogen

# Define a directory for logs
Environment="LOG_DIR=/home/ubuntu/PEAK-Assistant/autogen"

# Command to initialize pyenv, activate the virtual environment, and redirect stdout to a log file
ExecStart=/home/ubuntu/.pyenv/versions/peak-assistant/bin/python3 run_endpoint.py team-config-gpt-4.1.json --ssl-certfile cert.pem --ssl-keyfile key.pem --port 49822 --env-file ../.env

# Restart the service automatically if it crashes
Restart=always
RestartSec=5

# Environment variables (if needed)
Environment="PYTHONUNBUFFERED=1"

# User and group to run the serviceq
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
